# scripts/systemd/butler-b9-weekly-gate.timer
[Unit]
Description=Weekly B9 learning + SWE full LIVE when gate opens (Sunday 03:30)

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
