# scripts/systemd/butler-eval-sync.timer
[Unit]
Description=Weekly Butler WeChat dataset sync to LangFuse

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
