# scripts/systemd/butler-eval-sync.service
[Unit]
Description=Butler WeChat corpus → LangFuse dataset sync
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@WFXM_ROOT@
Environment=PYTHONPATH=@WFXM_ROOT@
EnvironmentFile=-@WFXM_ROOT@/.env
ExecStart=/bin/bash @WFXM_ROOT@/scripts/butler-wechat-dataset-sync.sh
StandardOutput=append:@WFXM_ROOT@/logs/butler-eval-sync.log
StandardError=append:@WFXM_ROOT@/logs/butler-eval-sync.log

[Install]
WantedBy=default.target
